# system/etc/profile.d/sudo-with-alias.sh
# allows using sudo with aliases
alias sudo='sudo '
