# system/etc/profile.d/sudo-with-alias.sh
# Allow sudo to work with aliases
alias sudo='sudo '
